# devel/etc/catkin/profile.d/cs4750.zsh
/home/bjh254/homework_ws/devel/.private/cs4750/etc/catkin/profile.d/cs4750.zsh